# deploy/prot.service
[Unit]
Description=prot voice conversation service
After=pipewire-pulse.service
Requires=pipewire-pulse.service

[Service]
Type=simple
WorkingDirectory=/home/cyan/workplace/prot
EnvironmentFile=/home/cyan/workplace/prot/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/cyan/workplace/prot/.venv/bin/uvicorn prot.app:app --host 0.0.0.0 --port 8000 --log-level info
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
